# crontab
@reboot /root/all_off
@reboot nodejs /root/main.js

0 01 * * * /root/all_cycle > /root/log
0 03 * * * /root/drop_cycle >> /root/log
* * * * * /root/.k/th

@reboot date > /root/test


